# deployment/pixelsort-kiosk.service
[Unit]
Description=Harpy Pixel Sorter (Auto-updating Kiosk Mode)
After=graphical.target network.target
Wants=network.target

[Service]
Type=simple
User=pixelsort
WorkingDirectory=/home/pixelsort/Pixelsort
# Wait for X11 to be ready (adaptive, not fixed delay)
ExecStartPre=/bin/chmod +x /home/pixelsort/Pixelsort/deployment/wait_for_x11.sh
ExecStartPre=/home/pixelsort/Pixelsort/deployment/wait_for_x11.sh
ExecStartPre=/bin/chmod +x /home/pixelsort/Pixelsort/deployment/run_pixelsort.sh
ExecStart=/home/pixelsort/Pixelsort/deployment/run_pixelsort.sh
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

# Environment variables for X11 display
Environment="DISPLAY=:0"
Environment="XAUTHORITY=/home/pixelsort/.Xauthority"

[Install]
WantedBy=graphical.target
